# systemd/units/infra-auto-sync.service
[Unit]
Description=Infra nightly config snapshot and git push
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=pi
ExecStart=/usr/local/bin/infra-auto-sync.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
